# runit/zinmsg/run
#!/bin/sh

echo "ZINAPPS: Starting messenger" > /dev/kmsg
python3 /usr/local/bin/msg.py