# runit/zinkey/run
#!/bin/sh

echo "ZINAPP: starting wireguard key manager" > /dev/kmsg

python3 /usr/local/bin/wgkey.py
